# supervisor/nginx.ini
[program:nginx]
process_name=%(program_name)s_%(process_num)02d
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
