# supervisor/swoole.ini
[program:swoole]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /app/artisan swoole:http start
user=app
autostart=true
autorestart=true
startretries=5
priority=15
stdout_events_enabled=true
stderr_events_enabled=true
